# launchagent/com.coretex.aos-organ.cortex.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.coretex.aos-organ.cortex</string>

    <key>ProgramArguments</key>
    <array>
        <string>/opt/homebrew/bin/node</string>
        <string>/Library/AI/AI-AOS/AOS-organ-dev/AOS-organ-cortex/AOS-organ-cortex-src/server/index.js</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Library/AI/AI-AOS/AOS-organ-dev/AOS-organ-cortex/AOS-organ-cortex-src</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>NODE_ENV</key>
        <string>development</string>
        <key>CORTEX_PORT</key>
        <string>4040</string>
    </dict>

    <key>StandardOutPath</key>
    <string>/tmp/com.coretex.aos-organ.cortex.out.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/com.coretex.aos-organ.cortex.err.log</string>

    <key>RunAtLoad</key>
    <false/>

    <key>KeepAlive</key>
    <false/>
</dict>
</plist>
